# systemd/sqlcoach.timer
[Unit]
Description=Run Daily SQL Coach

[Timer]
# Runs daily using the system timezone (UTC by default on EC2)
# 14:00 UTC is 08:00 CST / 09:00 CDT
OnCalendar=*-*-* 14:00:00
Persistent=true

[Install]
WantedBy=timers.target
